# systemd/ato-mcp-update.service
[Unit]
Description=ato-mcp: pull the latest ATO index from GitHub Releases
Documentation=https://github.com/gunba/ato-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/ato-mcp/environment
ExecStart=%h/.local/bin/ato-mcp update
ExecStartPost=/usr/bin/systemctl --user try-restart ato-mcp-serve.service
ExecStartPost=%h/.local/bin/ato-mcp stats
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
