# systemd/ato-mcp-serve.service
[Unit]
Description=ATO MCP server (HTTP daemon)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ato-mcp daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
