# systemd/ato-mcp-maintainer-weekly.timer
[Unit]
Description=ato-mcp: maintainer weekly incremental release

[Timer]
OnCalendar=Sun *-*-* 03:43:00
RandomizedDelaySec=30min
Persistent=true
Unit=ato-mcp-maintainer-weekly.service

[Install]
WantedBy=timers.target
