# systemd/legal-mcp-update.timer
[Unit]
Description=Australian Legal MCP: weekly corpus refresh (end user)

[Timer]
# Weekly at 08:17 local time, with randomized delay to avoid a stampede
# across installs that all run `legal-mcp update` at the same minute.
OnCalendar=Mon *-*-* 08:17:00
RandomizedDelaySec=30min
Persistent=true
Unit=legal-mcp-update.service

[Install]
WantedBy=timers.target
