# systemd/legal-mcp.service
[Unit]
Description=Australian Legal MCP private HTTP service
After=network-online.target local-fs.target
Wants=network-online.target
RequiresMountsFor=/var/lib/australian-legal-mcp
AssertPathExists=/var/lib/australian-legal-mcp/.legal-mcp-data-volume

[Service]
Type=simple
User=legal-mcp
Group=legal-mcp
EnvironmentFile=/etc/australian-legal-mcp/legal-mcp.env
ExecStart=/usr/local/bin/legal-mcp serve --port 51235 --require-ready-corpus
Restart=on-failure
RestartSec=5s
TimeoutStartSec=10min
TimeoutStopSec=45s
KillSignal=SIGTERM
UMask=0027
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
CapabilityBoundingSet=
SystemCallArchitectures=native
ReadOnlyPaths=/var/lib/australian-legal-mcp
ReadWritePaths=/var/lib/australian-legal-mcp/state
TasksMax=64
LimitNOFILE=8192
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
