# systemd/legal-mcp-maintainer-weekly.service
[Unit]
Description=Australian Legal MCP: weekly local refresh, build, and activation
Documentation=https://github.com/gunba/australian-legal-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=LEGAL_MCP_REPO_DIR=%h/src/australian-legal-mcp
Environment=LEGAL_MCP_PROJECT_DATA_DIR=%h/src/australian-legal-mcp/data
Environment=MALLOC_ARENA_MAX=24
ExecStart=%h/src/australian-legal-mcp/scripts/maintainer-sync.sh
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=43200
