# systemd/legal-mcp-maintainer-weekly.timer
[Unit]
Description=Australian Legal MCP: maintainer weekly incremental release

[Timer]
OnCalendar=Sun *-*-* 03:43:00
RandomizedDelaySec=30min
Persistent=true
Unit=legal-mcp-maintainer-weekly.service

[Install]
WantedBy=timers.target
