# systemd/legal-mcp-update.service
[Unit]
Description=Australian Legal MCP: pull the latest corpus from GitHub Releases
Documentation=https://github.com/gunba/australian-legal-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/australian-legal-mcp/environment
ExecStart=%h/.local/bin/legal-mcp update
ExecStartPost=/usr/bin/systemctl --user try-restart legal-mcp-serve.service
ExecStartPost=%h/.local/bin/legal-mcp stats
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
